# internal/deploy/assets/awg-agent.service
[Unit]
Description=AmneziaWG admin agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/awg-agent
EnvironmentFile=-/etc/awg-agent/awg-agent.env
Restart=on-failure
RestartSec=2s

# The agent manages WireGuard interfaces and needs elevated network
# capabilities, but does not need full root.
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/awg-agent /var/lib/awg-agent
StateDirectory=awg-agent
# The userspace agent (amneziawg-go) puts each interface's UAPI control socket
# under /var/run/amneziawg (a symlink to /run/amneziawg). ProtectSystem=strict
# makes /run read-only, so the agent can't create that dir itself — have systemd
# create it writable (and clean it up on stop). Harmless for the kernel agent,
# which drives netlink and never opens a UAPI socket.
RuntimeDirectory=amneziawg

[Install]
WantedBy=multi-user.target
